# rootfs/etc/s6-overlay/s6-rc.d/meticulous/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
exec python3 /usr/bin/run.py
